# scripts/passenger/upstart.conf
description "nginx http daemon"

start on runlevel [2]
stop on runlevel [016]

console output

exec /opt/nginx/sbin/nginx -c /opt/nginx/conf/nginx.conf -g "daemon off;"

respawn
